# root/etc/uci-defaults/30_luci-theme-vayu
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Vayu=/luci-static/vayu
	set luci.main.mediaurlbase=/luci-static/vayu
	commit luci
EOF
exit 0
